# openwrt/etc/init.d/vpn-control-boot
#!/bin/sh /etc/rc.common
# Renders the runtime xray config from persistent state at boot, before the
# xray service (START=91) tries to read it. /var/etc is tmpfs, so without this
# the office VPN would wait for the first cron tick after every reboot.

START=90
USE_PROCD=0

start() {
  [ -x /usr/libexec/vpnctl ] || return 0
  state="/etc/vpn-control/state.json"
  [ -s "$state" ] || return 0

  enabled="$(jq -r '.enabled' "$state" 2>/dev/null)"
  selected="$(jq -r '.selected // ""' "$state" 2>/dev/null)"
  [ "$enabled" = "true" ] && [ -n "$selected" ] || return 0

  . /usr/lib/vpn-control/common.sh
  . /usr/lib/vpn-control/renderer.sh
  ensure_dirs
  if render_config "$selected" "$(cfg '.socks_port' 10808)" "$XRAY_CONFIG"; then
    chmod 600 "$XRAY_CONFIG"
    logger -t vpnctl "boot: rendered config for $selected"
  else
    logger -t vpnctl "boot: failed to render config for $selected"
  fi
}
